# supervisor.conf
[program:djangovisor]
command=/virtualenv/path/bin/gunicorn -c gunicorn.conf.py wsgi:application
directory=/project/path
user=www-data
autostart=true
autorestart=true
redirect_stderr=true